# letsencrypt-lighttpd.service
[Unit]
Description=Letsencrypt manual renewal

[Service]
Type=oneshot
ExecStart=/home/user/letsencrypt-cloudflare/renew.sh
